# _site_setup/upstart_logicalc_live.conf
start on started postgresql
stop on stopped postgresql

respawn
respawn limit 3 5

setuid paul
setgid paul

exec /home/paul/live/_site_setup/service_demon_live.sh